# lock.service
[Unit]
Description=lock on suspend

[Service]
User=${USER}
ExecStart=/usr/bin/physlock -s

[Install]
WantedBy=sleep.target
